# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[inet_http_server]
port=127.0.0.1:9001

[supervisord]
logfile=./run/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=./run/pids/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:frontend]
command=pnpm run dev
directory=./frontend
user=root
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=./run/logs/frontend.stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=./run/logs/frontend.stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
environment=PATH="/root/.nvm/versions/node/v22.9.0/bin:/root/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
stopwaitsecs=15
stopsignal=TERM
killasgroup=true
stopasgroup=true

[program:agentserver]
command=uv run aetherterm --host localhost --port 57575 --unsecure --debug --more
directory=.
user=root
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=./run/logs/agentserver.stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=./run/logs/agentserver.stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
environment=PATH="/root/.nvm/versions/node/v22.9.0/bin:/root/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
stopwaitsecs=20
stopsignal=TERM
killasgroup=true
stopasgroup=true