# assets/aetherterm.service
[Unit]
Description=AetherTerm Terminal Server

[Service]
ExecStart=/usr/bin/aetherterm