# assets/butterfly.service
[Unit]
Description=Butterfly Terminal Server

[Service]
ExecStart=/usr/bin/aetherterm/server.py
